# roles/certbot/templates/certbot.service.j2
[Unit]
Description=Certbot Renewal Service
Documentation=https://certbot.eff.org/docs

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew --quiet --agree-tos

[Install]
WantedBy=multi-user.target